# src/servers/stargate/stargate-daemon.service
[Unit]
Description="STARGATE: Skynet AI Agent & Telemetry Hub"
After=network.target

[Service]
User=franklin
WorkingDirectory=/mnt/clusterfs2/workspace/gaming/game-skynet-minecraft
# TELEMETRY_PORT and CHONK_IP for telemetry filtering and listener binding
Environment="RCON_PASS=dinosaurExTraVaGanZa1969" "RCON_PORT=25575" "TELEMETRY_PORT=5005" "CHONK_IP=10.10.8.60"
# Ensure the project root is in PYTHONPATH
Environment="PYTHONPATH=/mnt/clusterfs2/workspace/gaming/game-skynet-minecraft"
ExecStart=/mnt/clusterfs2/workspace/gaming/game-skynet-minecraft/.venv/bin/python3 src/servers/stargate/stargate-daemon.py
Restart=always
RestartSec=10
StandardOutput=append:/mnt/clusterfs2/workspace/gaming/game-skynet-minecraft/logs/skynet_agent.log
StandardError=append:/mnt/clusterfs2/workspace/gaming/game-skynet-minecraft/logs/skynet_agent.log

[Install]
WantedBy=multi-user.target
