# src/servers/edge-t/edge-t-daemon.service
[Unit]
Description=Skynet Edge-T Vision & Construction Daemon (Asus Tinker)
After=network-online.target
Wants=network-online.target

[Service]
User=franklin
Group=engr

WorkingDirectory=/home/skynet
ExecStart=/home/skynet/.venv/bin/python3 -m src.schematics.npu_spatial_engine --query optimal_vector

StandardOutput=append:/home/skynet/logs/edge-t.log
StandardError=append:/home/skynet/logs/edge-t-error.log

Restart=always
RestartSec=10


# Resource management for the NXP i.MX 8M
CPUAccounting=true
MemoryAccounting=true

[Install]
WantedBy=multi-user.target
