# scripts/skynet-ai.service
[Unit]
Description=Skynet AI Orchestrator for Minecraft
After=network.target minecraft.service
Wants=minecraft.service

[Service]
Type=simple
User=franklin
Group=franklin
WorkingDirectory=/home/franklin/workspace/gaming/game-chonk-minecraft
# Use the absolute path to the virtual environment's python
ExecStart=/home/franklin/workspace/gaming/game-chonk-minecraft/venv/bin/python3 schematics/skynet_orchestrator.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skynet-ai

[Install]
WantedBy=multi-user.target
