# scripts/skynet-ai.service
[Unit]
Description=Skynet AI Process for Minecraft Orchestration
After=network.target

[Service]
Type=simple
User=franklin
Group=franklin
WorkingDirectory=/home/franklin/workspace/gaming/game-chonk-minecraft/schematics
ExecStart=/home/franklin/workspace/gaming/game-chonk-minecraft/venv/bin/python3 skynet_daemon.py
Restart=always
RestartSec=10
StandardOutput=append:/home/franklin/workspace/gaming/game-chonk-minecraft/logs/skynet_daemon.log
StandardError=append:/home/franklin/workspace/gaming/game-chonk-minecraft/logs/skynet_daemon.log
SyslogIdentifier=skynet-ai

[Install]
WantedBy=multi-user.target
